# com.luxafor.menu-daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.luxafor.menu-daemon</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/Users/larry.craddock/Projects/luxafor/luxafor-menu-daemon.sh</string>
    </array>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <true/>
    
    <key>StandardOutPath</key>
    <string>/tmp/luxafor-menu-daemon.log</string>
    
    <key>StandardErrorPath</key>
    <string>/tmp/luxafor-menu-daemon.log</string>
    
    <key>WorkingDirectory</key>
    <string>/Users/larry.craddock/Projects/luxafor</string>
</dict>
</plist>